# monitor.conf
[program:monitor]
command=python3 -u web.py
directory=/home/ubuntu/speedy
stdout_logfile=/home/ubuntu/logs/monitor.log
redirect_stderr=true
LOG_LEVEL=debug